# karoo.service
[Unit]
Description=Karoo Stratum V1 Proxy
After=network.target

[Service]
ExecStart=/usr/local/bin/karoo -config /etc/karoo/config.json -idle_grace_ms 15000
Restart=always
RestartSec=2
User=karoo
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target